# python/crontab
*/15 * * * *  parallel -a tmp/hosts  /usr/bin/python uptime_load_averages.py